# deployment/schedule-app.service
[Unit]
Description=ScheduleApp Service
After=network.target

[Service]
User=root
WorkingDirectory=/root/apps/schedule-app
ExecStart=/root/apps/schedule-app/scheduleApp
Environment="DB_HOST=db"
Environment="DB_PORT=5432"
Environment="DB_USER=schedule_user"
Environment="DB_PASSWORD=schedule_pass"
Environment="DB_NAME=schedule_db"
EnvironmentFile=-/root/apps/schedule-app/config.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target